# network/net-config@.service
[Unit]
Description=Scripts configuring network: %i

[Service]
Type=oneshot
ExecStart=/etc/net/%i
