# network/send-broadcast@.service
[Unit]
Description=Broadcast to %i

[Service]
Type=oneshot
ExecStart=ping -c 4 -b %i
