# network/set-wifi@.service
[Unit]
Description=Set WiFi: %i

[Service]
Type=oneshot
ExecStart=/usr/bin/wifi %i
