# network/flush-neigh/set-wifi-and-flush-neigh@.service
[Unit]
Description=Wifi %i and notify to remove neighbors

[Service]
Type=oneshot
ExecStart=/root/userconfig/network/flush-neigh/wifi-and-send.py %i
